# infra/systemd/tunescore-industry-scraper.timer
[Unit]
Description=Run TuneScore Industry Pulse scraper every 4 hours
Requires=tunescore-industry-scraper.service

[Timer]
# Run every 4 hours at minute 0
OnCalendar=*-*-* 00,04,08,12,16,20:00:00

# If the system was off during a scheduled run, run the job after boot
Persistent=true

# Randomize start time by up to 5 minutes to avoid API rate limit issues
RandomizedDelaySec=300

[Install]
WantedBy=timers.target

